# vps_dashboard.service
[Unit]
Description=VPS Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/VPS_Dashboard
ExecStart=/root/VPS_Dashboard/.venv/bin/python app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
